# Backend/crontab
*/30 * * * * cd /App && ./Backend >> /var/log/cron.log 2>&1
